# virtuoso/upstart-virtuoso.conf
# virtuoso - Upstart script for Virtuoso Open Source Edition 7.1

description "virtuoso"
author "Geir Ove Grønmo <grove...@bouvet.no>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Run as unprivileged user
# setuid vagrant
# setgid vagrant

# chdir /home/vagrant

# Start Virtuoso
exec /usr/local/virtuoso-opensource/bin/virtuoso-t -df +configfile /usr/local/virtuoso-opensource/var/lib/virtuoso/db/virtuoso.ini

pre-stop script
  exec echo "SHUTDOWN;" | /usr/local/virtuoso-opensource/bin/isql
end script
